# deploy/atlas-scheduler.timer
[Unit]
Description=ATLAS Scheduler Timer - Daily at 8:30 AM ET

[Timer]
OnCalendar=*-*-* 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
